# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:tor-socat]
command=ash /tor-socat.sh
